# deploy/systemd/wonderland-admin.service
[Unit]
Description=Wonderland Admin Panel (Next.js standalone)
After=network-online.target wonderland-api.service
Wants=network-online.target

[Service]
Type=simple
User=wonderland
Group=wonderland
WorkingDirectory=/srv/wonderland/admin
Environment=NODE_ENV=production
Environment=PORT=3003
Environment=HOSTNAME=127.0.0.1
Environment=NODE_OPTIONS=--max-old-space-size=320
EnvironmentFile=-/srv/wonderland/admin/.env.local
ExecStart=/opt/node/bin/node server.js

Restart=always
RestartSec=5

MemoryHigh=350M
MemoryMax=450M
OOMPolicy=restart

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
